# scripts/sync-backup.timer
[Unit]
Description=Run weekly backup sync

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
